# launchd/com.dicheeko.orchestra.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dicheeko.orchestra</string>

    <!-- Paths match this machine (`which node`); re-check if node moves -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shaughnyd/.local/node/bin/node</string>
        <string>--import</string>
        <string>tsx</string>
        <string>/Users/shaughnyd/Documents/GitHub/ORCHESTRA/src/index.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shaughnyd/Documents/GitHub/ORCHESTRA</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/shaughnyd/.local/node/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>ORCHESTRA_CONFIG</key>
        <string>/Users/shaughnyd/Documents/GitHub/ORCHESTRA/config/studio.yaml</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/shaughnyd/Documents/GitHub/ORCHESTRA/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shaughnyd/Documents/GitHub/ORCHESTRA/logs/launchd.err.log</string>
</dict>
</plist>
